# snapcitr-hotkey.service
[Unit]
Description=Snapcitr Global Hotkey Listener
After=graphical.target

[Service]
Type=simple
ExecStart=REPO_PATH/.venv/bin/python REPO_PATH/hotkey_listener.py
Restart=on-failure
Environment="DISPLAY=:0"
Environment="XAUTHORITY=%h/.Xauthority"

[Install]
WantedBy=default.target
